# mac-settings/crontab/switowski
# * * * * * command
# * - minute (0-59)
# * - hour (0-23)
# * - day of the month (1-31)
# * - month (1-12)
# * - day of the week (0-6, 0 is Sunday)
# command - command to execute

00 20 * * * osascript /Users/switowski/Library/Mobile\ Documents/com~apple~ScriptEditor2/Documents/Evernote\ Backup.scpt daily ~/Backups/Evernote/
